# contrib/sentinel.service
[Unit]
Description=Sentinel - terminal workspace
Documentation=https://github.com/opus-domini/sentinel
StartLimitIntervalSec=60
StartLimitBurst=4

[Service]
Type=simple
ExecStart=%h/.local/bin/sentinel
Restart=on-failure
RestartSec=2
KillMode=process
Environment=SENTINEL_LOG_LEVEL=info
Environment=HOME=%h
Environment=TERM=xterm-256color
Environment=LANG=C.UTF-8

# Hardening
# Note: ProtectHome, PrivateDevices and PrivateTmp are NOT set because
# Sentinel needs access to ~/.sentinel/, /dev/ptmx (PTY) and /tmp/tmux-* sockets.
# Note: NoNewPrivileges is intentionally omitted — sudo is required for multi-user
# session support (running tmux as a different OS user).
# Note: NoNewPrivileges and SystemCallArchitectures are intentionally
# omitted — sudo requires new-privilege capability for multi-user sessions.

[Install]
WantedBy=default.target
